# debian/raspistatusleds_0.1-1_arhf/DEBIAN/postinst
#!/bin/sh
update-rc.d raspiledsts defaults
update-rc.d readyled defaults
exit 0


