# debian/raspistatusleds_0.1-1_arhf/DEBIAN/postrm
#!/bin/sh
set -e
if [ "$1" = "remove" ]; then
	update-rc.d -f raspiledsts remove
	update-rc.d -f readyled remove
fi
exit 0
~                                      
